# debian/pvs-monitor.service
[Unit]
Description=PVS6 solar monitor
After=network-online.target
Wants=network-online.target

[Service]
User=pvs-monitoring
ExecStart=/usr/bin/pvs-monitor --config /etc/pvs-monitor/config.yaml --db /var/lib/pvs-monitor/readings.db
Restart=on-failure
RestartSec=30
RestartPreventExitStatus=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
