# debian/pvs-api.service
[Unit]
Description=PVS6 solar data API
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/pvs-api -db /var/lib/pvs-monitor/readings.db -tls-cert /var/lib/pvs-monitor/server.crt -tls-key /var/lib/pvs-monitor/server.key
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
